# config/systemd-creds.socket
[Socket]
ListenStream=/run/systemd-creds/socket
